# tg-keyword-monitor-bot.service
[Unit]
Description=tg-keyword-monitor-bot Service
After=network.target

[Service]
Type=simple
User=root
# 注意修改为你自己的目录
WorkingDirectory=/root/tg-keyword-monitor-bot
# 注意修改为你自己的目录
ExecStart=python3 /root/tg-keyword-monitor-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
